# bin/pironman5.service
[Unit]
Description=pironman5 service
After=multi-user.target

[Service]
Type=simple
User=pironman5
Group=pironman5
WorkingDirectory=/opt/pironman5
ExecStart=/usr/local/bin/pironman5 start
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
